# .gollum/container-gollum.service.example
# Example service
# Do not copy whole. Please compare to your own generated service.

# container-gollum.service
# autogenerated by Podman 4.5.1
# Sat Jun 17 00:20:42 +08 2023

[Unit]
Description=Podman container-gollum.service
Documentation=man:podman-generate-systemd(1)
RequiresMountsFor=%t/containers

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStart=/usr/bin/podman run \
	--cidfile=%t/%n.ctr-id \
	--cgroups=no-conmon \
	--rm \
	--sdnotify=conmon \
	-d \
	--replace \
	--name gollum \
	--security-opt label=disable \
	-v %h/wiki:/wiki \
	-v %h/wiki/.gollum:/etc/gollum \
	-p 4567:4567 localhost/gollum:v5.3.2-commonmark \
	--config /etc/gollum/config.rb
ExecStop=/usr/bin/podman stop \
	--ignore -t 10 \
	--cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm \
	-f \
	--ignore -t 10 \
	--cidfile=%t/%n.ctr-id
Type=notify
NotifyAccess=all
# Optional - see man systemd.resource-control
# Only allow local access
IPAddressDeny=any
IPAddressAllow=localhost
# RAM limit before throttling
MemoryHigh=128M

[Install]
WantedBy=default.target
